# addons/homepods/systemd/vhf-shipods.service
[Unit]
Description=ShiPods als OwnTone-Ausgang ausgewaehlt halten
After=owntone.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vhf-shipods.sh
