# addons/homepods/systemd/vhf-shipods.timer
[Unit]
Description=ShiPod-Auswahl regelmaessig pruefen

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
